# dockerfiles/vamp-pulse/supervisord.conf
[supervisord]
nodaemon=true

[program:vamp-pulse]
directory=/
command=java -Dstream.url="http://localhost:10001/v1/stats/stream" -jar pulse-server-assembly-0.7.7-rc1.jar > /dev/null 2>&1 &




