# clique-consul/supervisord.conf
[supervisord]
nodaemon=true

[program:Elasticsearch]
command=gosu elasticsearch elasticsearch
autorestart=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:Logstash]
command=gosu logstash logstash -f /etc/logstash/conf.d/logstash.conf
autorestart=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:Kibana]
command=kibana-start
autorestart=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:VampGatewayAgent]
command=/opt/vamp/vamp-gateway-agent --logo=false --storeType=consul --storeConnection=127.0.0.1:8500
autorestart=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:consul]
command=/bin/consul agent -server -config-dir=/config
autorestart=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0